# celeryd.conf
[program:celery]
command=celery worker -A app.celery --loglevel=INFO

numprocs=1
autostart=true
autorestart=true
